# cfg/usr/etc/profile.d/env.sh
export PYTHON_HISTORY='/dev/null'
export LC_ALL='C.UTF-8'
export EDITOR='mk'
